# deployment/plex-mounter.service
[Unit]
Description=Public Media Mounter Web Interface
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/findPublicMedia
Environment="PRODUCTION=true"
Environment="SECRET_KEY=CHANGE_THIS_TO_A_RANDOM_SECRET_KEY"
Environment="PORT=5001"
ExecStart=/usr/bin/python3 /opt/findPublicMedia/web_app_pro.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
